# ubo_app/system/services/ubo-system.service.tmpl
[Unit]
Description=Ubo System Service
Before=hardware-test.service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=root
Environment=INSTALLATION_PATH={{INSTALLATION_PATH}}
Environment=USERNAME={{USERNAME}}
ExecStart={{INSTALLATION_PATH}}/env/bin/ubo-system
WorkingDirectory={{INSTALLATION_PATH}}
KillMode=control-group
Restart=always
RestartPreventExitStatus=255
RuntimeDirectory=ubo
RuntimeDirectoryMode=0755
